# packaging/DEBIAN/postinst
#!/bin/bash
# ============================================================================
# postinst — Sentinel post-installation script.
#
# Executed by dpkg AFTER all package files have been copied to the system.
#
# Actions:
#   1. Create the quarantine directory with strict permissions.
#   2. Reload systemd to pick up the new sentinel.service unit.
#   3. Enable the sentinel daemon so it starts on boot.
#   4. Start the sentinel daemon immediately.
# ============================================================================

set -e

# ── 1. Create quarantine directory ──────────────────────────────────────────
# The quarantine directory stores copies of infected files.
# It must be owned by root with 700 permissions so that regular users
# cannot access quarantined malware samples.

QUARANTINE_DIR="/opt/quarantine"

if [ ! -d "$QUARANTINE_DIR" ]; then
    mkdir -p "$QUARANTINE_DIR"
    echo "[sentinel] Created quarantine directory: $QUARANTINE_DIR"
fi

chown root:root "$QUARANTINE_DIR"
chmod 700 "$QUARANTINE_DIR"
echo "[sentinel] Quarantine directory permissions set to 700 (root only)."

# ── 2. Reload systemd daemon ───────────────────────────────────────────────
# This tells systemd to re-scan /etc/systemd/system/ for new or changed
# unit files.  Required because we just dropped sentinel.service there.

systemctl daemon-reload
echo "[sentinel] systemd daemon reloaded."

# ── 3. Enable the sentinel service ────────────────────────────────────────
# Creates the symlink in multi-user.target.wants so the daemon starts
# automatically on every boot.

systemctl enable sentinel
echo "[sentinel] sentinel.service enabled (will start on boot)."

# ── 4. Start the daemon immediately ───────────────────────────────────────
# Start the daemon right now so the user doesn't have to reboot.
# We use '|| true' to prevent the postinst from failing if clamd is not
# yet fully initialised (the daemon handles this gracefully internally).

systemctl start sentinel || true
echo "[sentinel] sentinel.service started."

echo ""
echo "═══════════════════════════════════════════════════════════"
echo "  Sentinel installed successfully!"
echo ""
echo "  Daemon:  systemctl status sentinel"
echo "  GUI:     Launch 'Sentinel Security' from the app menu"
echo "           or run: /opt/sentinel-gui/sentinel-gui"
echo "═══════════════════════════════════════════════════════════"
echo ""

exit 0
