# packaging/DEBIAN/postinst
#!/bin/bash
# ============================================================================
# postinst — Post-Installation Script for Sentinel Endpoint Security
# ============================================================================
#
# Executed by dpkg AFTER all package files have been unpacked to their
# final locations.  This script:
#
#   1. Creates the quarantine directory with strict permissions (700).
#   2. Reloads systemd so it picks up the new sentinel.service unit.
#   3. Enables the sentinel service to start on boot.
#   4. Starts the sentinel daemon immediately.
#
# Exit on any error — dpkg will roll back the installation if this fails.
# ============================================================================

set -e

# ── 1. Create the quarantine vault ──────────────────────────────────────────
# The daemon stores quarantined malware here.  Only root should have access.
QUARANTINE_DIR="/opt/quarantine"

if [ ! -d "$QUARANTINE_DIR" ]; then
    echo "sentinel: Creating quarantine directory at $QUARANTINE_DIR ..."
    mkdir -p "$QUARANTINE_DIR"
fi

# Lock down permissions: only root can read/write/traverse.
# This prevents unprivileged users from accessing quarantined malware samples.
chmod 700 "$QUARANTINE_DIR"
chown root:root "$QUARANTINE_DIR"
echo "sentinel: Quarantine directory secured (mode 700, owner root:root)."

# ── 2. Reload systemd manager configuration ────────────────────────────────
# Required so systemd sees the newly installed sentinel.service unit file.
echo "sentinel: Reloading systemd daemon configuration ..."
systemctl daemon-reload

# ── 3. Enable the service to start on boot ──────────────────────────────────
# This creates a symlink in /etc/systemd/system/multi-user.target.wants/
echo "sentinel: Enabling sentinel.service for automatic start on boot ..."
systemctl enable sentinel.service

# ── 4. Start the daemon immediately ─────────────────────────────────────────
# So the user gets real-time protection right after installation.
echo "sentinel: Starting sentinel daemon ..."
systemctl start sentinel.service

echo ""
echo "╔══════════════════════════════════════════════════════════════════╗"
echo "║  Sentinel Endpoint Security installed successfully!             ║"
echo "║                                                                 ║"
echo "║  • Daemon running:  systemctl status sentinel                   ║"
echo "║  • Logs:            journalctl -u sentinel -f                   ║"
echo "║  • GUI:             Search 'Sentinel' in your app menu          ║"
echo "║  • Quarantine:      /opt/quarantine                             ║"
echo "╚══════════════════════════════════════════════════════════════════╝"
echo ""

exit 0
