# packaging/DEBIAN/prerm
#!/bin/bash
# ============================================================================
# prerm — Pre-Removal Script for Sentinel Endpoint Security
# ============================================================================
#
# Executed by dpkg BEFORE package files are removed from the filesystem.
# This script:
#
#   1. Stops the running sentinel daemon gracefully (SIGTERM).
#   2. Disables the service so it no longer starts on boot.
#
# We use "|| true" defensively — if the service was already stopped or
# never started, the uninstallation should still succeed.
#
# Exit on any error.
# ============================================================================

set -e

echo "sentinel: Stopping sentinel daemon ..."
systemctl stop sentinel.service || true

echo "sentinel: Disabling sentinel.service from boot ..."
systemctl disable sentinel.service || true

echo "sentinel: Reloading systemd daemon configuration ..."
systemctl daemon-reload || true

echo "sentinel: Pre-removal cleanup complete."

exit 0
